# systemd/.config/systemd/user/battery-osd.service
[Unit]
Description=Battery OSD Warning

[Service]
ExecStart=%h/.local/bin/battery-osd.sh
